# packages/wbsairback-config/wbsairback-config-1.1/etc/init.d/wbsairback-volumes-multipath
#! /bin/sh
#########################################
#.        .__--------___.         .     #
# \.     /             C \       / \.  /#
#   \   /           \__ ' \  .  /    \/ #
#    \_|  \   |-\   / \ \C'_/ \/        #
#------/ /-| /---| /-\ \----------------#
#      '_))'_))  '_)) '_))              #
#                                       #
#########################################
#        (c)WHITEBEARSOLUTIONS          #
#      GNU Generic Public License       #
#########################################
#
### BEGIN INIT INFO
# Provides:          wbsairback-volumes-multipath
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start WBSAirback volumes multipath daemon at boot time
# Description:       Enable WBSAirback volumes multipath.
### END INIT INFO
#
DAEMON_STATUS=0

if [ -f /etc/default/functions ]; then
  . /etc/default/functions
fi


start() {
	echo -n "Starting WBSAirback Volumes Multipath ["
	/usr/share/wbsairback/bin/volumes multipath noha > /dev/null 2>&1 || DAEMON_STATUS = 1
	service_wait 3
    echo -n "] "
    if [ $DAEMON_STATUS -eq 0 ]; then
        print_done
    else
        print_fail
    fi
}

stop() {
   	echo -n "Stopping WBSAirback Volumes Multipath ["
	service_wait 1
    echo -n "] "
    print_done
}

case "$1" in
  start)
    start
    ;;

  stop)
    stop
    ;;

 restart)
    stop
    sleep 1
    start
    ;;
  *)
    echo "Usage: /etc/init.d/wbsairback-volumes-multipath {start|stop|restart}"
    exit 1
    ;;
esac

exit 0


